# src/scripts/cccp_pkgupdatelistener.service
[Unit]
Description=cccp-pkgupdatelistener.service

[Service]
ExecStart=/opt/cccp-service/src/manage.py pkgupdatelistener
