# src/scripts/cccp_imagescanner.service
[Unit]
Description=cccp-imagescanner.service

[Service]
ExecStart=/opt/cccp-service/src/manage.py imagescanner
