# beanstalk_worker/cccp-build-worker.service
[Unit]
Description=cccp-build-worker.service

[Service]
ExecStart=/opt/cccp-service/beanstalk_worker/worker_start_build.py
