# beanstalk_worker/cccp-scan-worker.service
[Unit]
Description=cccp-scan-worker.service

[Service]
ExecStart=/opt/cccp-service/beanstalk_worker/worker_start_scan.py
