# beanstalk_worker/cccp-test-worker.service
[Unit]
Description=cccp-test-worker.service

[Service]
ExecStart=/opt/cccp-service/beanstalk_worker/worker_start_test.py
