# scripts/cccp_imagescanner.service
[Unit]
Description=cccp-imagescanner.service

[Service]
ExecStart=/opt/cccp-service/manage.py imagescanner
Restart=on-failure
