# scripts/cccp_triggerbuilds.service
[Unit]
Description=cccp-triggerbuilds.service

[Service]
ExecStart=/opt/cccp-service/manage.py triggerbuilds
Restart=on-failure
